# target/etc/init.d/mobmatter
#!/bin/sh /etc/rc.common

START=99 # after mobilus

USE_PROCD=1
PROG=/opt/jungi/bin/mobmatter
LOG=syslog

start_service() {
    procd_open_instance
    procd_set_param command $PROG
    procd_set_param env LOG=$LOG
    procd_close_instance
}
